# server_management/crudapp.service
# systemd service file, /etc/systemd/system/monolith.service
[Unit]
Description=Barebones crud app
After=network.target

[Service]
Type=simple
User=myuser
WorkingDirectory=/path/to/your/app
ExecStart=/path/to/your/app/executable
Restart=always
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
